# Apps/gandi/crontab
*/5 * * * * /home/pi/gandi/updateGANDI.sh
